# .docker/supervisor/app.ini
[program:app]
process_name=%(program_name)s_%(process_num)02d
command=php /home/www-data/artisan queue:work --verbose --tries %(ENV_QUEUE_TRIES)s --timeout %(ENV_QUEUE_TIMEOUT)s --sleep %(ENV_QUEUE_SLEEP_SECONDS)s --delay %(ENV_QUEUE_DELAY)s
autostart=true
autorestart=true
user=www-data
numprocs=%(ENV_SV_NUM_PROCS)s
redirect_stderr=true
stdout_logfile=/dev/stdout